# deploy/midinet-client.service
[Unit]
Description=MIDInet Client Daemon
After=network.target sound.target midinet-bridge.service
Wants=network.target
Requires=midinet-bridge.service

[Service]
Type=simple
User=midi
Group=midi
ExecStart=/usr/local/bin/midi-client --config /etc/midinet/client.toml
Restart=always
RestartSec=2

Environment=RUST_LOG=info

# Real-time scheduling
Nice=-20
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80
LimitRTPRIO=99
LimitMEMLOCK=infinity

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/midinet /tmp
NoNewPrivileges=false
PrivateTmp=false

StandardOutput=journal
StandardError=journal
SyslogIdentifier=midinet-client

[Install]
WantedBy=multi-user.target
