# deploy/midinet-update.path
[Unit]
Description=Watch for MIDInet update requests

[Path]
PathModified=/var/lib/midinet/update-trigger

[Install]
WantedBy=multi-user.target
